# assets/tab-atelier-headless.service
[Unit]
Description=tab-atelier headless tab daemon
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/tab-atelier-headless
Restart=on-failure
RestartSec=2s

User=tab-atelier
Group=tab-atelier

# Pin HOME into the state dir so platform::config_dir() / state_base_dir()
# land at /var/lib/tab-atelier/.{config,local/state} — under our own
# StateDirectory, never /root or /home.
Environment=HOME=/var/lib/tab-atelier
Environment=RUST_LOG=tab_atelier=info

# systemd creates /var/lib/tab-atelier, chowns it to the service user,
# and persists it across restarts (only wiped on dpkg `purge` via postrm).
StateDirectory=tab-atelier
StateDirectoryMode=0750
UMask=0077

# --- Hardening ---
# Privilege containment
NoNewPrivileges=true
CapabilityBoundingSet=
AmbientCapabilities=
RestrictSUIDSGID=true

# Filesystem — explicit bind-mount sandbox.
#
# TemporaryFileSystem=/ replaces the root with a tmpfs so nothing
# from the host is visible by default. Then BindReadOnlyPaths
# selectively grants read-only access to the system paths the
# daemon and spawned shells need (binaries, libraries, /etc for
# passwd / resolv.conf / SSL trust store). BindPaths selectively
# grants read-write access to the user's actual work dirs.
#
# To open a tab in a new directory, add that dir to BindPaths via
# `systemctl edit tab-atelier-headless` (or extend this list and
# restart). The daemon's own state survives under StateDirectory
# (mounted automatically into the namespace).
TemporaryFileSystem=/
BindReadOnlyPaths=/usr /lib /lib64 /bin /sbin /etc /opt
# Optional but commonly needed at run time:
#   /run/systemd/resolve  — symlink target for /etc/resolv.conf on
#                            systemd-resolved systems; without it
#                            DNS lookups from the shell fail.
#   /sys/fs/cgroup        — read-only cgroup info (some tools probe
#                            it; harmless when absent).
# Leading `-` makes the bind optional, so the unit still starts on
# systems where the path doesn't exist.
BindReadOnlyPaths=-/run/systemd/resolve -/sys/fs/cgroup
# Read-write work directories. Add yours here (or via
# `systemctl edit`). /srv is the conventional spot for site data;
# extend with /home, /opt/projects, etc. as you go.
#
# The daemon's own home — /var/lib/tab-atelier — is writable
# automatically: StateDirectory= mounts it into the namespace
# regardless of TemporaryFileSystem=/, so tabs that `cd $HOME` or
# write to `$HOME/...` (the tab-atelier service user's home) work.
BindPaths=/srv
PrivateTmp=true
# NOTE: PrivateDevices=true would mask /dev/ptmx + /dev/pts and break
# every PTY this daemon opens. Leave it off and rely on the rest of
# the protections for the kernel/device attack surface.
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
ProtectClock=true
ProtectHostname=true
ProtectProc=invisible
ProcSubset=pid

# Network — TCP only (happier-relay is loopback, API binds INADDR_ANY).
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6

# Namespaces / IPC / keyring
RestrictNamespaces=true
RestrictRealtime=true
RemoveIPC=true
KeyringMode=private
LockPersonality=true
MemoryDenyWriteExecute=true
SystemCallArchitectures=native

# Allowlist with the standard service-syscall baseline, then carve
# out the dangerous groups. @system-service already includes the
# openpty / fork / exec / mmap calls this daemon needs.
SystemCallFilter=@system-service
SystemCallFilter=~@clock @cpu-emulation @debug @module @mount @obsolete @privileged @raw-io @reboot @resources @swap @keyring @pkey
SystemCallErrorNumber=EPERM

[Install]
WantedBy=multi-user.target
