# assets/tab-atelier-headless.service
[Unit]
Description=tab-atelier headless tab daemon
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/tab-atelier-headless
Restart=on-failure
RestartSec=2s

User=tab-atelier
Group=tab-atelier

# Pin HOME into the state dir so platform::config_dir() / state_base_dir()
# land at /var/lib/tab-atelier/.{config,local/state} — under our own
# StateDirectory, never /root or /home.
Environment=HOME=/var/lib/tab-atelier
Environment=RUST_LOG=tab_atelier=info

# systemd creates /var/lib/tab-atelier, chowns it to the service user,
# and persists it across restarts (only wiped on dpkg `purge` via postrm).
StateDirectory=tab-atelier
StateDirectoryMode=0750
UMask=0077

# --- Hardening ---
# Privilege containment
NoNewPrivileges=true
# CAP_NET_ADMIN lets the daemon program the per-tab nftables egress
# allowlist (see src/net_nft.rs). It's the ONLY capability granted. Tabs do
# NOT keep it: each tab's shell is spawned through `setpriv --ambient-caps=
# -all` so an agent can't touch the firewall. Leave both empty to disable
# egress-allowlist enforcement entirely (tabs still spawn, just unconfined).
CapabilityBoundingSet=CAP_NET_ADMIN
AmbientCapabilities=CAP_NET_ADMIN
RestrictSUIDSGID=true

# Filesystem — explicit bind-mount sandbox.
#
# TemporaryFileSystem=/ replaces the root with a tmpfs so nothing
# from the host is visible by default. Then BindReadOnlyPaths
# selectively grants read-only access to the system paths the
# daemon and spawned shells need (binaries, libraries, /etc for
# passwd / resolv.conf / SSL trust store). BindPaths selectively
# grants read-write access to the user's actual work dirs.
#
# To open a tab in a new directory, add that dir to BindPaths via
# `systemctl edit tab-atelier-headless` (or extend this list and
# restart). The daemon's own state survives under StateDirectory
# (mounted automatically into the namespace).
TemporaryFileSystem=/
BindReadOnlyPaths=/usr /lib /lib64 /bin /sbin /etc /opt
# Optional but commonly needed at run time:
#   /run/systemd/resolve  — symlink target for /etc/resolv.conf on
#                            systemd-resolved systems; without it
#                            DNS lookups from the shell fail.
#   /sys/fs/cgroup        — read-only cgroup info (some tools probe
#                            it; harmless when absent).
# Leading `-` makes the bind optional, so the unit still starts on
# systems where the path doesn't exist.
BindReadOnlyPaths=-/run/systemd/resolve -/sys/fs/cgroup

# Local sockets users typically expect to reach from a tab — IPC
# with daemons running on the host. Both paths covered for docker
# because /var/run is usually a symlink to /run, but with our
# private tmpfs the symlink doesn't follow until we bind the
# target into the namespace too. Extend with the same pattern if
# you use podman, libvirt, etc.
BindPaths=-/run/docker.sock -/var/run/docker.sock
# Read-write work directories. Add yours here (or via
# `systemctl edit`). /srv is the conventional spot for site data;
# extend with /home, /opt/projects, etc. as you go.
#
# The daemon's own home — /var/lib/tab-atelier — is writable
# automatically: StateDirectory= mounts it into the namespace
# regardless of TemporaryFileSystem=/, so tabs that `cd $HOME` or
# write to `$HOME/...` (the tab-atelier service user's home) work.
BindPaths=/srv
PrivateTmp=true
# NOTE: PrivateDevices=true would mask /dev/ptmx + /dev/pts and break
# every PTY this daemon opens. Leave it off and rely on the rest of
# the protections for the kernel/device attack surface.
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
ProtectClock=true
ProtectHostname=true
ProtectProc=invisible
ProcSubset=pid

# Network — AF_INET/AF_INET6 for the API, AF_UNIX for IPC, AF_NETLINK so
# the daemon can program the per-tab nftables egress allowlist (nft talks
# to the kernel over NETLINK_NETFILTER). The tabs themselves only need
# INET/UNIX and hold no caps, so they can't do anything privileged over
# netlink even though the family is permitted service-wide.
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6 AF_NETLINK

# Per-tab resource limits (cgroup v2). Delegate the service's own
# cgroup subtree so the daemon can create one child cgroup per tab and
# write memory.max / cpu.max / pids.max into it (see src/cgroup.rs).
# Delegation grants write access to *only* this unit's subtree and
# coexists with ProtectControlGroups= above (systemd keeps the rest of
# /sys/fs/cgroup read-only). Without this the per-tab limit settings are
# silently inert — the daemon detects the non-writable subtree and skips
# limiting, so the unit still works if you remove this line.
Delegate=yes

# Namespaces / IPC / keyring
RestrictNamespaces=true
RestrictRealtime=true
RemoveIPC=true
KeyringMode=private
LockPersonality=true
MemoryDenyWriteExecute=true
SystemCallArchitectures=native

# Allowlist with the standard service-syscall baseline, then carve
# out the dangerous groups. @system-service already includes the
# openpty / fork / exec / mmap calls this daemon needs.
SystemCallFilter=@system-service
SystemCallFilter=~@clock @cpu-emulation @debug @module @mount @obsolete @privileged @raw-io @reboot @resources @swap @keyring @pkey
# Re-allow capset + the privilege-DROP syscalls (setuid/setgid/...). They're
# in @privileged (dropped above), but common interactive tools call them at
# startup just to SHED privileges — `ping` (iputils), `sudo`, `su`, and our
# own `setpriv` cap-strip — so blocking the syscall makes them error out
# ("ping: setuid: Operation not permitted") even when they're not trying to
# gain anything. Re-allowing the syscalls grants nothing: raising to another
# uid still needs CAP_SETUID (tabs hold none) and NoNewPrivileges blocks
# setuid-binary escalation. ICMP itself already works via unprivileged ping
# sockets (net.ipv4.ping_group_range), so this is purely about the drop.
SystemCallFilter=capset setuid setgid setgroups setresuid setresgid setreuid setregid setfsuid setfsgid
SystemCallErrorNumber=EPERM

[Install]
WantedBy=multi-user.target
